# setup/servidor/uptime-kuma-webhook.service
[Unit]
Description=Uptime Kuma -> Telegram webhook (THDORA)
After=network.target docker.service

[Service]
Type=simple
User=varopc
WorkingDirectory=/home/varopc/yggdrasil-dew
EnvironmentFile=/home/varopc/yggdrasil-dew/.env
ExecStart=/usr/bin/python3 /home/varopc/yggdrasil-dew/scripts/uptime-kuma-webhook.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
